# basic-topology/volumes/supervisord.conf
[supervisord]
nodaemon=true

[program:zebra]
command=/usr/sbin/zebra -f /etc/quagga/zebra.conf